# runner/app/runit/infer/run
#!/bin/sh

echo "Starting infer"
exec 2>&1
cd /app
exec python -u -m app.live.infer